# crumbd@.service
[Unit]
Description=crumb daemon

[Service]
ExecStart=/usr/local/bin/crumbd %I

[Install]
WantedBy=multi-user.target
