# crumb@.service
[Unit]
Description=Crumb Daemon

[Service]
ExecStart=/usr/local/bin/crumbd %I
Restart=always

[Install]
WantedBy=multi-user.target
