# crumb.service
[Unit]
Description=Crumb Daemon

[Service]
ExecStart=/home/dobyrch/Projects/crumb/crumbd
Restart=always

[Install]
WantedBy=multi-user.target
